# packaging/postinstall.sh
systemctl daemon-reload
systemctl enable tart-guest-agent.service
systemctl start tart-guest-agent.service
