# packaging/preremove.sh
systemctl stop tart-guest-agent.service
systemctl disable tart-guest-agent.service
